# systemd/zdt-api.service
[Unit]
Description=ZDT API Server v1.0.0 - Standalone API for ZDT Mobile
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/zdt-project/zdt-api
ExecStart=/usr/local/bin/gunicorn --bind 0.0.0.0:2000 --workers 2 --timeout 120 server:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
