# systemd/zdt-api.service
[Unit]
Description=ZDT API Server - Standalone API for ZDT Mobile
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__ZDT_API_DIR__
# Multi-worker with Gunicorn config file (gunicorn_config.py)
# - SQLite WAL mode supports concurrent reads; writes serialized via DB lock
# - Redis rate limiter (optional) for shared state across workers
# - Preloaded app for faster startup
ExecStart=__PYTHON__ -m gunicorn --config gunicorn_config.py server:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
